# compose/fullroute/supervisord.conf
[supervisord]
nodaemon=true

[program:gobgpd]
command=/usr/bin/start-gobgpd
